# default-files/etc/init.d/freemesh
#!/bin/sh /etc/rc.common

START=99 
STOP=99 

start() {
        /usr/share/freemesh/Startup.sh
}

stop() {
        echo stop
        # commands to kill application
}