# default-files/etc/uci-defaults/00-www-perms.sh
#!/bin/sh 

chmod -R go+rX /www

exit 0
